# config/auto-unlock.service
[Unit]
Description=Automated Lock Screen Recovery Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/glitch/.config/i3/scripts/auto-unlock.sh
Restart=always
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target